# ansible/roles/sentry/templates/sentry.ini
[program:sentry-web]
directory=/var/www/sentry/
command=/var/www/sentry/bin/sentry start
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile syslog
stderr_logfile syslog

[program:sentry-worker]
directory=/var/www/sentry/
command=/var/www/sentry/bin/sentry celery worker -B
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile syslog
stderr_logfile syslog
